# systemd/nginx-ab-editor.service
[Unit]
Description=nginx antiblocking editor

[Service]
Type=simple
ExecStart=nginx-ab-editor-remote

[Install]
WantedBy=multi-user.target
